# backup/dot-config/systemd/user/vps-backup.service
[Unit]
Description=User backup job

[Service]
Type=oneshot
EnvironmentFile=%h/.config/restic/env
ExecStart=%h/opt/scripts/backup-now.bash
